# autostart/rosita.desktop
[Desktop Entry]
Type=Application
Name=Rosita
Exec=/bin/bash /home/pi/Nextcloud/rosita/rosita.sh
X-GNOME-Autostart-enabled=true
